# examples/autoreload-nginx/proxy/runit/services/nginx/run
#!/bin/sh
source /etc/envvars;
nginx -g 'daemon off;'