# examples/applications-as-runit-services/docker/services/__app__/run
#!/bin/sh

cd /usr/node/$APP
# # wait for core service
# sv up core || exit 1
exec node app.js