# examples/monorepo-container-hmr-services/runit/lib/run
#!/bin/sh

cd /usr/node/$APP;
source /etc/envvars;
exec yarn run lib:dev;